# geoip-update.timer
[Unit]
Description=GeoIP Update Service timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
